# deploy/supervisor-config.conf
# Supervisor配置文件 - AI助理系统守护进程
# 文件位置：/etc/supervisor/conf.d/ai-assistant.conf

[program:ai-assistant]
# 进程名称
process_name=ai-assistant

# 启动命令（使用绝对路径）
command=/usr/bin/python3 /var/www/ai-assistant/assistant_web.py

# 工作目录
directory=/var/www/ai-assistant

# 运行用户（建议使用www-data或nobody，不要用root）
user=www-data

# 自动启动
autostart=true

# 自动重启（如果进程崩溃）
autorestart=true

# 启动重试次数
startretries=3

# 停止信号
stopsignal=INT

# 停止等待时间（秒）
stopwaitsecs=10

# 重定向stderr到stdout
redirect_stderr=true

# 日志文件
stdout_logfile=/var/log/ai-assistant.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10

stderr_logfile=/var/log/ai-assistant-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10

# 环境变量（如果需要）
environment=
    PYTHONUNBUFFERED="1",
    LANG="zh_CN.UTF-8",
    LC_ALL="zh_CN.UTF-8"

# 进程优先级
priority=999

# 备注：
# 1. 确保Python3已安装：python3 --version
# 2. 确保依赖已安装：pip3 install pymysql
# 3. 确保MySQL数据库已配置
# 4. 确保uploads目录权限正确：chown -R www-data:www-data /var/www/ai-assistant/uploads
